# fedora/src/home/config/autostart/gnome-keyring-daemon.desktop
[Desktop Entry]
Type=Application
Name=gnome-keyring-daemon
Exec=/dockerstartup/gnome-keyring.sh