# latest/src/home/config/autostart/disable-screensaver.desktop
[Desktop Entry]
Type=Application
Name=disable-screensaver
Comment=Disable screensaver
Terminal=true
Exec=xset s off
